# crontab
@hourly /usr/local/bin/php /usr/src/ipcheck/reporter.php
